# setup/resources/crontab
# vim:syntax=crontab

# vim:syntax=crontab

# Delete all files and directories in the ~/.temporary dir
0 */3 * * * /bin/rm -rf ~/.temporary/*

# Delete files in trash monthly
@monthly /bin/rm -rf ~/.Trash/*

0 */24 * * * cd ~/.dotfiles/setup && /usr/local/bin/brew bundle dump --force --describe --quiet
