# setup/resources/crontab
# vim:syntax=crontab

PATH="/usr/local/bin"

# Delete all files and directories in the ~/.temporary dir
0 */12 * * * /bin/rm -rf ~/.temporary/*
# Delete files in trash monthly
@monthly /bin/rm -r ~/.Trash/*
0 */24 * * * cd ~/.dotfiles/setup && brew bundle dump --force
